# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

exec \
	sh /etc/openvpn/startvpn.sh


